# llm-proxy.service
[Unit]
Description=LLM traffic logging reverse proxy
After=network.target docker.service 
Wants=docker.service

[Service]
Type=simple

# --- Adjust these paths to match your installation ---
ExecStart=%h/opt/llm-opensearch-logger/venv/bin/python %h/opt/llm-opensearch-logger/llm_proxy.py
WorkingDirectory=%h/opt/llm-opensearch-logger

# Environment overrides (edit as needed)
# For multiple backends, use PROXY_MAPPINGS instead of PROXY_PORT/OLLAMA_PORT
# Format: "listen_port:upstream_port[:name]" or "listen_port:host:upstream_port[:name]"
Environment=PROXY_MAPPINGS=11434:11435:ollama,8002:8003:llama-server
Environment=OPENSEARCH_URL=http://localhost:9200
Environment=OPENSEARCH_INDEX=ollama-traffic
Environment=LOG_BODY_MAX_BYTES=65536

Restart=on-failure
RestartSec=5

# Give Ollama time to start before the proxy attempts connections
ExecStartPre=/bin/sleep 3

StandardOutput=journal
StandardError=journal
SyslogIdentifier=ollama-proxy

[Install]
WantedBy=default.target
